# config/Docker/all-in-one-snapshot-image/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisor.log
# loglevel=DEBUG

[program:httpd]
command=/usr/sbin/apache2ctl -DFOREGROUND

[program:sshd]
command=/usr/sbin/sshd -D

[program:SecurityTokenService]
command=/usr/bin/java -DIAM_MODE=PROD -Dhazelcast.config=hazelcast.xml -DIAM_CONFIG=securitytokenservice.PROD.properties -jar /home/SecurityTokenService/SecurityTokenService.jar
user=SecurityTokenService
directory=/home/SecurityTokenService

[program:UserAdminService]
command=/usr/bin/java -DIAM_MODE=PROD -DCONSTRETTO_TAGS=DEV -jar /home/UserAdminService/UserAdminService.jar
user=UserAdminService   
directory=/home/UserAdminService

[program:UserAdminWebApp]
command=/usr/bin/java -DIAM_MODE=PROD -DIAM_CONFIG=/home/UserAdminWebApp/useradminwebapp.PROD.properties -jar /home/UserAdminWebApp/UserAdminWebApp.jar
user=UserAdminWebApp
directory=/home/UserAdminWebApp

[program:UserIdentityBackend]
command=/usr/bin/java -DCONSTRETTO_TAGS=DEV -jar UserIdentityBackend.jar 
user=UserIdentityBackend
directory=/home/UserIdentityBackend

[program:CRMService]
command=/usr/bin/java -DCONSTRETTO_TAGS=DEV -DIAM_MODE=PROD -jar Whydah-CRMService.jar
user=CRMService
directory=/home/CRMService

[program:StatisticsService]
command=/usr/bin/java -DCONSTRETTO_TAGS=DEV -jar Whydah-StatisticsService.jar
user=StatisticsService
directory=/home/StatisticsService

[program:TestWebApp]
command=/usr/bin/java -DIAM_MODE=TEST -DIAM_CONFIG=/home/TestWebApp/testwebapp.TEST.properties -jar /home/TestWebApp/TestWebApp.jar
user=TestWebApp
directory=/home/TestWebApp

[program:SSOLoginWebApp]
command=/usr/bin/java -DIAM_MODE=PROD -DIAM_CONFIG=/home/SSOLoginWebApp/ssologinservice.PROD.properties -jar /home/SSOLoginWebApp/SSOLoginWebApp.jar
user=SSOLoginWebApp
directory=/home/SSOLoginWebApp

[program:PostgreSQL]
command=/usr/lib/postgresql/10/bin/postgres -D /var/lib/postgresql/10/main --config-file=/etc/postgresql/10/main/postgresql.conf
user=postgres

[program:SPAProxyService]
command=/usr/bin/java -jar /home/SPAProxyService/Whydah-SPAProxyService.jar
user=SPAProxyService
directory=/home/SPAProxyService