# scripts/dev-test-server-scripts-automation-example/_etc_init_/StatisticsService.conf
# Ubuntu upstart file, http://upstart.ubuntu.com/cookbook/
description "StatisticsServer Upstart Configuration File"
author "baardl <bard.lind@gmail.com>"

start on runlevel [2345]
stop on runlevel [06]

script
    su -c "nohup /usr/lib/jvm/java-1.8.0-openjdk.x86_64/bin/java  -DIAM_MODE=PRO
D -DCONSTRETTO_TAGS=DEV -jar /home/StatisticsService/Whydah-StatisticsService.ja
r" -s /bin/sh - StatisticsService
end script
